# ansible/roles/app/templates/hippo-pwd.service
[Unit]
Description=uWSGI instance to serve Hippo IDP
After=network.target

[Service]
User={{ type }}service
Group=www-data
WorkingDirectory=/home/{{ type }}service/password-self-service/password_reset_frontend
Environment="PATH=/home/{{ type }}service/password-self-service/password_reset_frontend/bin"
ExecStart=/home/{{ type }}service/password-self-service/password_reset_frontend/bin/uwsgi --ini web.ini

[Install]
WantedBy=multi-user.target

